# cups/rootfs/etc/services.d/cups-mdns/run
#!/usr/bin/with-contenv bash
set -e

exec /usr/local/bin/publish-cups-mdns.sh
