# cups/rootfs/etc/services.d/avahi/run
#!/usr/bin/with-contenv bash
set -e

mkdir -p /run/avahi-daemon /etc/avahi/services

exec /usr/sbin/avahi-daemon --no-drop-root --no-chroot --no-rlimits --no-dbus
